# dist/smtptunnel-server.service
[Unit]
Description=SMTP Tunnel Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/smtptunnel-server run -c /etc/smtptunnel/config.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
